# setup/reader_egress2.service
[Unit]
Description=Site Sentinel Box Reader - Egress2
Requires=network.target

[Service]
SyslogIdentifier=reader_egress2
Type=forking
User=ubuntu
PIDFile=/home/ubuntu/site-sentinel-box-usb-readers/pids/reader_egress2.pid
WorkingDirectory=/home/ubuntu/site-sentinel-box-usb-readers
ExecStart=/bin/ruby /home/ubuntu/site-sentinel-box-usb-readers/reader_egress2_control.rb start
ExecReload=/bin/ruby /home/ubuntu/site-sentinel-box-usb-readers/reader_egress2_control.rb restart
ExecStop=/bin/ruby /home/ubuntu/site-sentinel-box-usb-readers/reader_egress2_control.rb stop
Restart=always

[Install]
WantedBy=multi-user.target